# deploy/systemd/pulse-agent.service
[Unit]
Description=Pulse Slack agent (Socket Mode)
After=network-online.target pulse-mcp.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/pulse/pulse-agent
EnvironmentFile=/opt/pulse/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
