# deploy/systemd/pulse-mcp.service
[Unit]
Description=Pulse MCP server (Neon-backed tools)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/pulse
EnvironmentFile=/opt/pulse/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run mcp:start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
